# deployments/docker/build/s6files/services.d/dsiem/run
#!/usr/bin/with-contenv sh
exec /dsiem/dsiem serve

